# sanitized-copypaste.service
[Unit]
Description=SanitizedCopyPaste Hotkey Daemon
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=peterdew
Group=peterdew
WorkingDirectory=/home/peterdew/repos/SanitizedCopyPaste-clean
ExecStart=/usr/bin/python3 /home/peterdew/repos/SanitizedCopyPaste-clean/hotkey_daemon.py
Restart=always
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
